# capacitr/templates/default/supervisor.conf.erb
[program:<%= @user %>]
command=/home/<%= @user %>/venv/bin/gunicorn_django -b 0.0.0.0:<%= @port %> -w 8
directory=/home/<%= @user %>/site
environment=PATH="/home/<%= @user %>/venv/bin"
user=<%= @user %>
autostart=true
autorestart=true
startsecs=5
exitcodes=0,2
numprocs=1

