# capacitr/templates/default/supervisor.local.conf.erb
[program:<%= node[:username] %>]
command=/home/<%= node[:username] %>/venv/bin/gunicorn_django -b 0.0.0.0:8000 -w 5
directory=/home/<%= node[:username] %>/site
environment=PATH="/home/<%= node[:username] %>/venv/bin"
user=<%= node[:username] %>
autostart=true
autorestart=true
startsecs=1
exitcodes=0,2
numprocs=1

