# capacitr/templates/default/supervisor.local.conf.erb
[program:<%= @username %>]
command=/home/<%= @username %>/venv/bin/gunicorn_django -b 0.0.0.0:8000 -w 5
directory=/home/<%= @username %>/site
environment=PATH="/home/<%= @username %>/venv/bin"
user=<%= @username %>
autostart=true
autorestart=true
startsecs=1
exitcodes=0,2
numprocs=1

